# services/temperature-controller-restarter@.path
[Path]
PathChanged=/etc/controller-setpoints/
PathChanged=/etc/%I

[Install]
WantedBy=multi-user.target
