# services/temperature-controller-restarter.path
[Path] 
PathChanged=/etc/setpoint
PathChanged=/etc/controller.conf

[Install] 
WantedBy=multi-user.target
